# ops/systemd/exports-cleanup.timer
[Unit]
Description=Clean up exports daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=osm-export-tool.target
