# ops/systemd/exports_cleanup.timer
[Unit]
Description=Run cleanup every minute

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
