# ops/systemd/exports-cleanup.service
[Unit]
Description=Clean up after exports

[Service]
Type=oneshot
ExecStart=/usr/bin/find /mnt-storage/stage -type d -ctime +7 -delete
